# rootfs/etc/services.d/litestream/run
#!/command/with-contenv bash
# Litestream SQLite replication service

# Skip if ENABLE_SPACES is not true
if [ "${ENABLE_SPACES:-false}" != "true" ]; then
  echo "[litestream] Disabled (set ENABLE_SPACES=true for persistence)"
  s6-svc -Od .
  exit 0
fi

# Skip if persistence not configured
if [ -z "$LITESTREAM_ACCESS_KEY_ID" ] || [ -z "$SPACES_BUCKET" ]; then
  echo "[litestream] ENABLE_SPACES=true but LITESTREAM_ACCESS_KEY_ID or SPACES_BUCKET not set"
  s6-svc -Od .
  exit 0
fi

echo "[litestream] Starting SQLite replication..."

# Run litestream replicate (continuous replication)
exec s6-setuidgid moltbot litestream replicate -config /etc/litestream.yml
