# rootfs/etc/services.d/moltbot/run
#!/command/with-contenv bash
# Moltbot gateway service

echo "[moltbot] Starting moltbot gateway..."

# Run as moltbot user with proper environment
export HOME=/home/moltbot
exec s6-setuidgid moltbot bash -c '
  source /home/moltbot/.bashrc
  source /home/moltbot/.profile
  export GRADIENT_API_KEY="${GRADIENT_API_KEY:-ChangeMe}"
  exec moltbot gateway --allow-unconfigured
'
