# rootfs/etc/services.d/moltbot/run
#!/command/with-contenv bash
# Moltbot gateway service

echo "[moltbot] Starting moltbot gateway..."

# Wait briefly for Tailscale to initialize
sleep 2

# Run gateway directly (litestream runs as separate service)
exec s6-setuidgid moltbot moltbot gateway --allow-unconfigured
