# rootfs/etc/services.d/crond/run
#!/command/execlineb -P
# Start cron daemon in foreground
/usr/sbin/cron -f
