# rootfs/etc/services.d/crond/run
#!/command/with-contenv bash
# Start cron daemon in foreground
exec /usr/sbin/cron -f > /proc/1/fd/1 2>&1
