# rootfs/etc/services.d/backup/run
#!/command/with-contenv bash
# Periodic backup service

# Skip if persistence not configured
if [ -z "$LITESTREAM_ACCESS_KEY_ID" ] || [ -z "$SPACES_BUCKET" ]; then
  echo "[backup] Persistence not configured, backup service disabled"
  s6-svc -Od .
  exit 0
fi

BACKUP_INTERVAL=${BACKUP_INTERVAL_SECONDS:-300}
echo "[backup] Starting backup loop (every ${BACKUP_INTERVAL}s)..."

# Main backup loop
while true; do
  sleep "$BACKUP_INTERVAL"
  /usr/local/bin/moltbot-backup
done
