# rootfs/etc/services.d/tailscale/run
#!/command/with-contenv bash
# Tailscale daemon (required for networking)

echo "[tailscale] Starting Tailscale daemon..."
exec /usr/local/bin/containerboot
