# rootfs/etc/services.d/tailscale/run
#!/command/with-contenv bash
# Tailscale daemon (required for networking)

echo "[tailscale] Starting Tailscale daemon..."

exec s6-setuidgid moltbot /usr/local/bin/containerboot
